# scriptman.service
[Unit]
Description=ScriptMan - Shell Script Manager
After=network.target

[Service]
Type=simple
User=www-data                          # change to your user
WorkingDirectory=/opt/scriptman        # directory with scriptman.py + scripts.json
ExecStart=/usr/bin/python3 scriptman.py
Restart=on-failure
RestartSec=5s

Environment=SCRIPTMAN_HOST=127.0.0.1
Environment=SCRIPTMAN_PORT=7000
Environment=SCRIPTMAN_CONFIG=/opt/scriptman/scripts.json
# Generate once and hardcode: python3 -c "import secrets; print(secrets.token_hex(32))"
Environment=SCRIPTMAN_SECRET=REPLACE_WITH_RANDOM_HEX

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
